# rootfs/etc/services.d/mariadb/run
#!/usr/bin/with-contenv sh

exec /scripts/wrapper/mysqld.sh
